# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
pidfile=/run/mt5-proxy/supervisord.pid
childlogdir=/logs

[program:xvfb]
priority=10
command=/usr/local/bin/run-xvfb.sh
autostart=true
autorestart=true
startsecs=3
startretries=999
stopsignal=TERM
stdout_logfile=/logs/xvfb.supervisor.log
stderr_logfile=/logs/xvfb.supervisor.err

[program:fluxbox]
priority=20
command=/usr/local/bin/run-fluxbox.sh
autostart=true
autorestart=true
startsecs=3
startretries=999
stopsignal=TERM
stdout_logfile=/logs/fluxbox.supervisor.log
stderr_logfile=/logs/fluxbox.supervisor.err

[program:x11vnc]
priority=30
command=/usr/local/bin/run-x11vnc.sh
autostart=true
autorestart=true
startsecs=3
startretries=999
stopsignal=TERM
stdout_logfile=/logs/x11vnc.supervisor.log
stderr_logfile=/logs/x11vnc.supervisor.err

[program:novnc]
priority=40
command=/usr/local/bin/run-novnc.sh
autostart=true
autorestart=true
startsecs=3
startretries=999
stopsignal=TERM
stdout_logfile=/logs/novnc.supervisor.log
stderr_logfile=/logs/novnc.supervisor.err

[program:api]
priority=50
command=/usr/local/bin/run-api.sh
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=999
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/api.supervisor.log
stderr_logfile=/logs/api.supervisor.err

[program:mt5-stack]
priority=60
command=/usr/local/bin/mt5-stack.sh
autostart=true
autorestart=true
startsecs=10
startretries=999
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/mt5-stack.supervisor.log
stderr_logfile=/logs/mt5-stack.supervisor.err
