# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/logs/supervisord.log
pidfile=/run/mt5-proxy/supervisord.pid
childlogdir=/logs
loglevel=info

[unix_http_server]
file=/run/mt5-proxy/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/mt5-proxy/supervisor.sock

[program:xvfb]
command=/usr/local/bin/run-xvfb.sh
user=trader
priority=10
autostart=true
autorestart=true
startsecs=2
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/xvfb.log
stderr_logfile=/logs/xvfb.err.log

[program:fluxbox]
command=/usr/local/bin/run-fluxbox.sh
user=trader
priority=20
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/fluxbox.log
stderr_logfile=/logs/fluxbox.err.log

[program:x11vnc]
command=/usr/local/bin/run-x11vnc.sh
user=trader
priority=30
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/x11vnc.log
stderr_logfile=/logs/x11vnc.err.log

[program:novnc]
command=/usr/local/bin/run-novnc.sh
user=trader
priority=40
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/novnc.log
stderr_logfile=/logs/novnc.err.log

[program:api]
command=/opt/mt5-proxy-venv/bin/uvicorn mt5_proxy.main:app --host 0.0.0.0 --port 8000 --proxy-headers
user=trader
directory=/app
priority=50
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/api.log
stderr_logfile=/logs/api.err.log

[program:mt5stack]
command=/usr/local/bin/mt5-stack.sh
user=trader
priority=60
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/mt5-stack.log
stderr_logfile=/logs/mt5-stack.err.log

[program:mt5bridge]
command=/usr/local/bin/mt5-bridge-watch.sh
user=trader
priority=70
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/mt5-bridge-watch.supervisor.log
stderr_logfile=/logs/mt5-bridge-watch.err.log

[program:healthwatch]
command=/usr/local/bin/healthwatch.sh
priority=100
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/logs/healthwatch.log
stderr_logfile=/logs/healthwatch.err.log
